# debian/python-oauth.upstart
description	"python-oauth"
author		"Danny Meloy <danny.meloy@bbc.co.uk>"

start on runlevel [2345]
stop on runlevel [016]
respawn

script
	exec sudo /usr/bin/nmosoauth
end script
